# deploy/pixel.service
[Unit]
Description=Bildfilter Stempel-Wanderung (pixel.rafeldt.ch)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pixel/stempel-wanderung
EnvironmentFile=/etc/pixel.env
ExecStart=/opt/pixel/.venv/bin/gunicorn \
  --workers 3 \
  --timeout 180 \
  --bind 127.0.0.1:5001 \
  --access-logfile - \
  --error-logfile - \
  app:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
